# supervisord.conf
[program:mkn_octane]
command=php /home/ergane/mneme_kai_nous/artisan octane:start
directory=/home/ergane/mneme_kai_nous/
user=ergane
autostart=true
autorestart=true
